# solar-telegram-listener.service
[Unit]
Description=Solar miner Telegram listener (two-way status queries)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/solar_miner_controller
ExecStart=/usr/bin/python3 /home/admin/solar_miner_controller/telegram_listener.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
